# supervisor.conf
[supervisord]
nodaemon=true

[program:app]
command=npm run up
